# install/systemd/feeder-daily-summary.service.in
[Unit]
Description=Send daily ADS-B feeder summary via Gotify
# managed-by-feeder-dashboard

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 @REPO_ROOT@/dashboard/daily-summary.py
Environment=FEEDER_REPO_ROOT=@REPO_ROOT@
EnvironmentFile=-@REPO_ROOT@/feeder.env
