# install/systemd/feeder-dashboard.timer.in
[Unit]
Description=Refresh feeder dashboard status every 30s
# managed-by-feeder-dashboard

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=feeder-dashboard-update.service

[Install]
WantedBy=timers.target
